# init_container_service/init_container.service
[Unit]
Description=docker

[Service]
ExecStart=/user/sbin/init_container.sh

[Install]
WantedBy=multi-user.target